# ufrgs-autostart.cron
35 * * * * root /usr/local/bin/ufrgs-autostart
